# .config/systemd/user/workenv.target
[Unit]
Description=Work Environment Target
Requires=workenv.service private.service

[Install]
WantedBy=default.target
